# agent/backup-agent.service
[Unit]
Description=MySQL Backup Agent
After=network.target

[Service]
Type=simple
# 若通过 Makefile 编译为了单执行文件，假设放在 /usr/local/bin 或某个部署目录
# 请根据实际情况修改下方的执行路径
ExecStart=/usr/local/bin/backup-agent

# 如果使用 .env 配置文件，可以指定其所在目录（可选）
# WorkingDirectory=/opt/autoBackup/agent

# 如果您不想使用 .env 文件或命令行参数，也可以直接在这里写入环境变量
# Environment="API_BASE=http://192.168.1.100:8000"
# Environment="TOKEN=您的安全密钥"
# Environment="HOSTNAME=bj-mysql-01"

# 自动重启与日志
Restart=always
RestartSec=5
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=backup-agent

[Install]
WantedBy=multi-user.target
